# com.erikjosephson.georgiapracticechecker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.erikjosephson.georgiapracticechecker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/erikjosephson/georgia-practice-checker/georgia_practice_checker.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run every Monday at 9:00 AM -->
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/erikjosephson/georgia-practice-checker/georgia_checker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/erikjosephson/georgia-practice-checker/georgia_checker.log</string>
</dict>
</plist>
